# src/program-redis.conf
;; -*- ini -*-
[program:redis]
user = nobody
command = /usr/sbin/redis-server /etc/redis.conf
autostart = true
autorestart = true
